# files/system/usr/lib/systemd/system/atomik-fish-setup.service
[Unit]
Description=Atomik OS - Fish shell setup per utenti esistenti
After=local-fs.target systemd-user-sessions.service

[Service]
Type=oneshot
ExecStart=/usr/bin/atomik-fish-setup
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
